# uxplay.service
[Unit]
Description=uxplay

[Service]
ExecStart=usr/bin/uxplay -n Raum311 -nh

[Install]
WantedBy=multi-user.target
